# ansible-supervisor-amzn2/templates/etc/supervisor/conf_d/welcome.conf.j2
[program:welcome]
command=echo "hi muthu"