# ansible-supervisor-amzn2/templates/etc/supervisor/conf_d/welcome.conf.j2
[program:welcome]
command=python -m SimpleHTTPServer 8080
autostart=yes
autorestart=yes
startretries=3